# deploy/systemd/aw-planner.service
[Unit]
Description=AgentWorks backlog planner pass (plan_work.sh) — needs planner.enabled=true

[Service]
Type=oneshot
# TODO(adopter): set these two to your runner user and clone path
User=agentworks
WorkingDirectory=/opt/agentworks/repo
EnvironmentFile=/etc/agentworks/aw.env
ExecStart=/opt/agentworks/repo/scripts/plan_work.sh
