# deploy/systemd/aw-triage.service
[Unit]
Description=AgentWorks agent-triage loop (triage_work.sh) — needs auto_triage.agent_triage=true
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# TODO(adopter): set these two to your runner user and clone path
User=agentworks
WorkingDirectory=/opt/agentworks/repo
EnvironmentFile=/etc/agentworks/aw.env
ExecStart=/opt/agentworks/repo/scripts/triage_work.sh
Restart=on-failure
RestartSec=60
MemoryMax=2G

[Install]
WantedBy=multi-user.target
